# dev/systemd/solar-protocol.service.d/api.conf
# /etc/systemd/service/solar-protocol.service.d/api.conf
# LOGLEVEL=DEBUG|INFO|WARNING|ERROR|CRITICAL

[Service]
Environment="LOGLEVEL=Warning"
